# data/packetfence/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
redirect_stderr=true
stdout_logfile=/var/log/mysql/supervisor.log
stderr_logfile=/var/log/mysql/supervisor.error_log
autorestart=false
environment=HOME="/etc/mysql",USER="mysql"
umask=007

[program:packetfence_start]
command=service packetfence start
startsecs=0
autorestart=false

[program:packetfence_stop]
command=service packetfence stop
startsecs=0
autostart=false
autorestart=false
